# tools/BootScripts/Resources/Defaults/Settings/inittab
id:2:initdefault:

l1:S:wait:/System/Links/Executables/BootDriver BootUp
su:S:respawn:/sbin/sulogin

l2:12345:wait:/System/Links/Executables/BootDriver Console
l6:6:wait:/System/Links/Executables/BootDriver Reboot
l0:0:wait:/System/Links/Executables/BootDriver Halt

ca:12345:ctrlaltdel:/sbin/shutdown -t1 -r now

1:2345:respawn:/System/Links/Executables/agetty tty1 9600
2:2345:respawn:/System/Links/Executables/agetty tty2 9600
3:2345:respawn:/System/Links/Executables/agetty tty3 9600
4:2345:respawn:/System/Links/Executables/agetty tty4 9600
5:2345:respawn:/System/Links/Executables/agetty tty5 9600
6:2345:respawn:/System/Links/Executables/agetty tty6 9600