# tools/TinyScripts/Resources/Defaults/Settings/inittab
#
# Boot-time system configuration/initialization script.
# This is run first except when booting in single-user mode.
#
::sysinit:/System/Settings/BootScripts/BootUp

# /bin/sh invocations on selected ttys
::askfirst:-/System/Links/Executables/sh
tty1::respawn:-/System/Links/Executables/getty tty1 9600 linux
tty2::respawn:-/System/Links/Executables/getty tty2 9600 linux
tty3::respawn:-/System/Links/Executables/getty tty3 9600 linux
tty4::respawn:-/System/Links/Executables/getty tty4 9600 linux
tty5::respawn:-/System/Links/Executables/getty tty5 9600 linux
tty6::respawn:-/System/Links/Executables/getty tty6 9600 linux

# Stuff to do when restarting the init process
::restart:/System/Links/Executables/init

# Stuff to do before rebooting
::ctrlaltdel:/System/Settings/BootScripts/Reboot
::shutdown:/System/Settings/BootScripts/Halt
